# mantis/configs/mantis.service
[Unit]
Description=Mantis DNS Sinkhole
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=mantis
Group=mantis
ExecStart=/usr/local/bin/mantis --config /etc/mantis/config.toml --data-dir /var/lib/mantis
Restart=always
RestartSec=5
LimitNOFILE=65535
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
